# system/supervisor/amcp.conf
[program:amcp-rpi]
command=/home/pi/amcp-rpi/server.py
autostart=true
autorestart=true
redirect_stderr=true
directory=/home/pi/amcp-rpi


